# roles/infestor/templates/infestor.service.j2
[Unit]
Description=Run ansible-pull regularly

[Service]
Type=oneshot
ExecStart=ansible-pull -U {{ infestor_giturl }} {{ infestor_playbook }}
