# roles/infestor/templates/infestor.timer.j2
[Unit]
Description=Timer to run infestor.service every 20 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec={{ infestor_interval }}
Unit=infestor.service

[Install]
WantedBy=timers.target
